# systemd/trxd.service
[Unit]
Description=trxd
After=network.target

[Service]
ExecStart=/usr/sbin/trxd
Type=forking
Restart=always

[Install]
WantedBy=default.target
