# debian/trx-control.prerm
#!/bin/sh

systemctl disable --now trxd
exit 0
