# debian/trx-control.prerm
#!/bin/sh

systemctl disable --now trx-control
exit 0
